# src/deb/control/postinst
systemctl daemon-reload
systemctl enable jjmpc.service
systemctl start jjmpc.service