# src/deb/control/prerm
systemctl stop jjmpc.service
systemctl disable jjmpc.service